# deploy/dgx2-worker.service
[Unit]
Description=DGX2 Worker - AEON vLLM file analysis consumer
After=network.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStartPre=%h/src/consensus-mcp/deploy/bootstrap-dirs.sh
ExecStart=/usr/bin/python3 %h/src/consensus-mcp/dgx2_worker.py
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3
RuntimeMaxSec=3600
StandardOutput=journal
StandardError=journal
Environment=COORDINATOR_URL=http://192.168.50.1:9000
Environment=AEON_URL=http://localhost:8001/v1
Environment=WORK_INBOX=%h/work-inbox
Environment=WORK_OUTBOX=%h/work-outbox
Environment=WORK_ARCHIVE=%h/work-archive
Environment=LEASE_DURATION=600
Environment=HEARTBEAT_INTERVAL=60
Environment=MAX_BATCH=100

[Install]
WantedBy=default.target
